# etc/profile.d/win-kex.sh
#!/bin/bash
export XAUTHORITY=~/.Xauthority
